# systemd/niri-mapper.user.service
# niri-mapper user service
#
# Installation (non-NixOS):
#   mkdir -p ~/.config/systemd/user
#   cp niri-mapper.user.service ~/.config/systemd/user/niri-mapper.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now niri-mapper
#
# Prerequisites:
#   - User must be in the 'input' group: sudo usermod -aG input $USER
#   - Log out and back in after group change
#   - niri-mapperd binary must be in /usr/bin (or adjust ExecStart path)
#   - Config file at ~/.config/niri-mapper/config.kdl
#
# Management:
#   systemctl --user status niri-mapper
#   systemctl --user restart niri-mapper
#   journalctl --user -u niri-mapper -f

[Unit]
Description=niri-mapper input remapping daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/niri-mapperd --config %h/.config/niri-mapper/config.kdl --foreground
Restart=on-failure
RestartSec=5

# Security hardening (user service compatible)
PrivateTmp=true

[Install]
WantedBy=graphical-session.target
