# systemd/niri-mapper.user.service
[Unit]
Description=niri-mapper input remapping daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/niri-mapperd --config %h/.config/niri-mapper/config.kdl --foreground
Restart=on-failure
RestartSec=5

# Security hardening (user service compatible)
PrivateTmp=true

[Install]
WantedBy=graphical-session.target
